# templates/maildev.ini.j2
[program:maildev]
command=/usr/bin/maildev -w {{ mailcatcher_port }} --ip 0.0.0.0 -s 25
autostart=true