# templates/maildev.ini.j2
[program:maildev]
command=/usr/bin/maildev -w {{ dork_ports['maildev'] }} --ip 0.0.0.0 -s {{ dork_ports['smtp'] }}
autostart=true