# templates/maildev.ini.j2
[program:maildev]
command=/usr/bin/maildev -w {{ mailcatcher_port }} --ip 127.0.0.1
autostart=true